# playbooks/dev.timer
[Unit]
Description="Run script every minute"

[Timer]
OnUnitActiveSec=1min

[Install]
WantedBy=default.target